# appendix/automake.conf
description "automake for /home/ueda/work"
author  "Ryuichi UEDA <ueda-ryuichi@aiit.ac.jp>"

start on runlevel [2345]
stop on runlevel [016]

chdir /
exec /home/ueda/work/automake
respawn
